# scripts/quantclaw.service
[Unit]
Description=QuantClaw AI Hardware Agent
Documentation=https://github.com/quantclaw/quantclaw
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
SupplementaryGroups=gpio spi i2c
WorkingDirectory=/home/pi/quantclaw
ExecStart=/usr/local/bin/quantclaw gateway --host 0.0.0.0 --port 42617
Restart=on-failure
RestartSec=5
EnvironmentFile=/home/pi/quantclaw/.env
Environment=RUST_LOG=info

# Expand ~ in config path
Environment=HOME=/home/pi

[Install]
WantedBy=multi-user.target
