# scripts/quantclaw.service
[Unit]
Description=QuantClaw AI Hardware Agent
Documentation=https://github.com/quantclaw/quantclaw
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__RPI_USER__
SupplementaryGroups=gpio spi i2c
WorkingDirectory=__RPI_DIR__
ExecStart=__RPI_DIR__/quantclaw gateway --host __GATEWAY_HOST__ --port __GATEWAY_PORT__
Restart=on-failure
RestartSec=5
EnvironmentFile=__RPI_DIR__/.env
Environment=RUST_LOG=info

# HOME is injected by the install/render script to match the selected SSH user.
Environment=HOME=__RPI_HOME__

[Install]
WantedBy=multi-user.target
